# install/com.yap.dictation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent for yap (macOS).
     Install:
       cp install/com.yap.dictation.plist ~/Library/LaunchAgents/
       # edit the yap path below if `which yap` differs
       launchctl load ~/Library/LaunchAgents/com.yap.dictation.plist
     Requires Microphone + Accessibility permission for the launching process. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yap.dictation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/yap</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/yap.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/yap.log</string>
</dict>
</plist>
